# overlay/etc/systemd/system/NetworkManager.service.d/nm-service-override.conf
[Service]
ExecStartPre=/bin/bash /usr/local/sbin/nm-wlanap-pre.sh
ExecStartPre=/bin/bash /usr/local/sbin/nm-lte-pre.sh
